# scripts/images/build-profiles/ubuntu/rootfs/etc/systemd/system/thin-resolv.service
# Writes /etc/resolv.conf from networkd lease state.
# systemd-resolved is masked in our images so thin-resolv always runs.
[Unit]
Description=Generate resolv.conf from networkd lease state
After=systemd-networkd.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c "sed -n 's/^DNS=\\(.*\\)/\\1/p' /run/systemd/netif/state | tr ' ' '\\n' | sed '/^$/d; s/^/nameserver /' > /etc/resolv.conf"
