# scripts/images/build-profiles/debian/_networkd/rootfs/etc/systemd/system/thin-resolv.path
[Unit]
Description=Watch for networkd state changes to update resolv.conf

[Path]
PathChanged=/run/systemd/netif/state
Unit=thin-resolv.service

[Install]
WantedBy=multi-user.target
